# server/charon.rubin.ru/etc/cron.hourly/drupal
#!/bin/bash

/usr/bin/wget -O - -q http://localhost/cron.php
