# server/webber.rubin.ru/etc/config-archive/etc/init.d/rmnologin
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-src/rc-scripts/init.d/rmnologin,v 1.9.4.1 2005/01/30 22:17:37 vapier Exp $

depend() {
	need localmount
}

start() {
	if [[ -f /etc/nologin.boot ]] ; then
		rm -f /etc/nologin /etc/nologin.boot &> /dev/null
	fi
}


# vim:ts=4
